# overlay/etc/systemd/system/axp221-monitor.service
[Unit]
Description=AXP221 Power Button Monitor
After=sysinit.target axp221-configure-pek.service
# I2C bus number is dynamic — script handles detection internally
Wants=systemd-modules-load.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/axp221-monitor.sh
# Prevent restart loops if IRQ is stuck
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=300
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
